# autosleep_allow
#!/sbin/openrc-run

depend()
{
	after sysfsconf
}

start()
{
	ebegin "Changing permissions on /sys/power/{autosleep,wake_(un)lock} to allow id=10000"
	chown 10000:10000 /sys/power/autosleep
	chown 10000:10000 /sys/power/wake_lock
	chown 10000:10000 /sys/power/wake_unlock
	eend $?
}

stop()
{
	ebegin "Nothing to stop, this service is just for changing permissions on /sys/power/{autosleep,wake_(un)lock}"
	eend $?
}

restart()
{
	ebegin "Nothing to restart, this service is just for changing permissions on /sys/power/{autosleep,wake_(un)lock}"
	eend $?
}
